# .config/autostart/flatpak-fcitx5.desktop
[Desktop Entry]
Type=Application
Name=fcitx5
Exec=flatpak run org.fcitx.Fcitx5
Terminal=false
